# launchd/com.lifedashboard.fatsecret.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
    DEPRECATED. Этот plist оставлен в репозитории как backup и для
    локальной отладки (`launchctl load ...`), но в продакшене FS-sync
    запускается из GitHub Actions — см. .github/workflows/fatsecret-sync.yml,
    cron 22:00 МСК. На локальном Mac этот job НЕ нужен и НЕ должен быть
    активирован, иначе данные будут писаться дважды.

    Если когда-то захочется снова поднять локально:
        cp launchd/com.lifedashboard.fatsecret.plist ~/Library/LaunchAgents/
        launchctl load ~/Library/LaunchAgents/com.lifedashboard.fatsecret.plist
    Снять с локального запуска (если стоит):
        launchctl unload ~/Library/LaunchAgents/com.lifedashboard.fatsecret.plist
        rm ~/Library/LaunchAgents/com.lifedashboard.fatsecret.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lifedashboard.fatsecret</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3.11</string>
        <string>/Users/admin/Desktop/LifeDashboard/scripts/fatsecret_scraper.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/admin/Desktop/LifeDashboard</string>

    <!-- 22:00 локального времени — рацион к этому моменту обычно записан. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>22</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/admin/Desktop/LifeDashboard/data/logs/fatsecret.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/admin/Desktop/LifeDashboard/data/logs/fatsecret.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
